# systemd/joystick-notify-tray.service
# ~/.config/systemd/user/joystick-notify-tray.service
[Unit]
Description=Joystick notify tray icon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/joystick-notify-tray
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target






